# dot_config/systemd/user/watchniriconfig.service
[Unit]
Description=Run script when niri files change

[Service]
ExecStart=/home/brennandunn/.local/bin/niri-configs-changed
